# api/systemd/kalembang.service
[Unit]
Description=Kalembang Alarm Controller
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kalembang/api
Environment="PATH=/opt/kalembang/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Safety: Force GPIO off before starting (optional extra safety)
ExecStartPre=/bin/sh -c 'gpio mode 4 out 2>/dev/null; gpio write 4 0 2>/dev/null; gpio mode 5 out 2>/dev/null; gpio write 5 0 2>/dev/null || true'

# Main service
ExecStart=/opt/kalembang/api/.venv/bin/uvicorn kalembang.main:app --host 0.0.0.0 --port 8088

# Restart policy
Restart=always
RestartSec=5

# Stop signals
TimeoutStopSec=10
KillMode=mixed

[Install]
WantedBy=multi-user.target
